# conf/supervisord.conf
[supervisord]
childlogdir = /home/work/open-falcon/logs/
logfile = /home/work/open-falcon/logs/supervisord.log
pidfile = /var/run/supervisord.pid
nodaemon = true
directory = /home/work/open-falcon

[unix_http_server]
file = /var/run/supervisor.sock
chmod = 0777
chown= root:root

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl = unix:///var/run/supervisor.sock

[program:start]
command = bash /home/work/scripts/start.sh
autorestart = false
priority = 1


[program:redis]
user = redis
command = /usr/bin/redis-server /etc/redis.conf
stdout_logfile = /home/work/open-falcon/logs/%(program_name)s.log
stderr_logfile = /home/work/open-falcon/logs/%(program_name)s.err
autostart=false
priority = 100
autorestart = true

[program:alarm]
user = root
directory = /home/work/open-falcon/%(program_name)s
command = /home/work/open-falcon/%(program_name)s/falcon-%(program_name)s -c /home/work/open-falcon/%(program_name)s/cfg.json
stdout_logfile = /home/work/open-falcon/logs/%(program_name)s.log
stderr_logfile = /home/work/open-falcon/logs/%(program_name)s.err
autostart=false
autorestart = true

[program:fe]
user = root
directory = /home/work/open-falcon/%(program_name)s
command = /home/work/open-falcon/%(program_name)s/falcon-%(program_name)s -c /home/work/open-falcon/%(program_name)s/cfg.json
stdout_logfile = /home/work/open-falcon/logs/%(program_name)s.log
stderr_logfile = /home/work/open-falcon/logs/%(program_name)s.err
autostart=false
autorestart = true

[program:graph]
user = root
directory = /home/work/open-falcon/%(program_name)s
command = /home/work/open-falcon/%(program_name)s/falcon-%(program_name)s -c /home/work/open-falcon/%(program_name)s/cfg.json
stdout_logfile = /home/work/open-falcon/logs/%(program_name)s.log
stderr_logfile = /home/work/open-falcon/logs/%(program_name)s.err
autostart=false
autorestart = true

[program:hbs]
user = root
directory = /home/work/open-falcon/%(program_name)s
command = /home/work/open-falcon/%(program_name)s/falcon-%(program_name)s -c /home/work/open-falcon/%(program_name)s/cfg.json
stdout_logfile = /home/work/open-falcon/logs/%(program_name)s.log
stderr_logfile = /home/work/open-falcon/logs/%(program_name)s.err
autostart=false
autorestart = true

[program:judge]
user = root
directory = /home/work/open-falcon/%(program_name)s
command = /home/work/open-falcon/%(program_name)s/falcon-%(program_name)s -c /home/work/open-falcon/%(program_name)s/cfg.json
stdout_logfile = /home/work/open-falcon/logs/%(program_name)s.log
stderr_logfile = /home/work/open-falcon/logs/%(program_name)s.err
autostart=false
autorestart = true

[program:query]
user = root
directory = /home/work/open-falcon/%(program_name)s
command = /home/work/open-falcon/%(program_name)s/falcon-%(program_name)s -c /home/work/open-falcon/%(program_name)s/cfg.json
stdout_logfile = /home/work/open-falcon/logs/%(program_name)s.log
stderr_logfile = /home/work/open-falcon/logs/%(program_name)s.err
autostart=false
autorestart = true

[program:sender]
user = root
directory = /home/work/open-falcon/%(program_name)s
command = /home/work/open-falcon/%(program_name)s/falcon-%(program_name)s -c /home/work/open-falcon/%(program_name)s/cfg.json
stdout_logfile = /home/work/open-falcon/logs/%(program_name)s.log
stderr_logfile = /home/work/open-falcon/logs/%(program_name)s.err
autostart=false
autorestart = true

[program:transfer]
user = root
directory = /home/work/open-falcon/%(program_name)s
command = /home/work/open-falcon/%(program_name)s/falcon-%(program_name)s -c /home/work/open-falcon/%(program_name)s/cfg.json
stdout_logfile = /home/work/open-falcon/logs/%(program_name)s.log
stderr_logfile = /home/work/open-falcon/logs/%(program_name)s.err
autostart=false
autorestart = true

[program:dashboard]
user = root
directory = /home/work/open-falcon/%(program_name)s
command = /home/work/open-falcon/%(program_name)s/env/bin/gunicorn -c gunicorn.conf wsgi:app -t 6000 
stdout_logfile = /home/work/open-falcon/logs/%(program_name)s.log
stderr_logfile = /home/work/open-falcon/logs/%(program_name)s.err
autostart=false
autorestart = true

[program:portal]
user = root
directory = /home/work/open-falcon/%(program_name)s
command = /home/work/open-falcon/%(program_name)s/env/bin/gunicorn -c gunicorn.conf wsgi:app -t 6000 
stdout_logfile = /home/work/open-falcon/logs/%(program_name)s.log
stderr_logfile = /home/work/open-falcon/logs/%(program_name)s.err
autostart=false
autorestart = true


[program:mailsender]
user = root
directory = /home/work/open-falcon/%(program_name)s
command = /home/work/open-falcon/%(program_name)s/falcon-%(program_name)s -c /home/work/open-falcon/%(program_name)s/cfg.json
stdout_logfile = /home/work/open-falcon/logs/%(program_name)s.log
stderr_logfile = /home/work/open-falcon/logs/%(program_name)s.err
autostart=false
autorestart = true
